# .config/upstart/kill-emacs-server.conf
description "Logout Kill Emacs Server Task"
start on session-end
task
script
  /usr/bin/emacsclient.emacs24 --eval "(kill-emacs)"
end script
